# tests/deployment/cliqz-backend.conf
description "{{ service_name }} test server"
start on runlevel [2345]
stop on shutdown

pre-start script
    while [ ! -f "/vagrant/tests/mockserver/server.py" ]
    do
        echo "Waiting for /vagrant share..."
        /bin/sleep 2s
    done

    echo "Share /vagrant available"
end script

exec python /vagrant/tests/mockserver/server.py {{ port }}
respawn